# .config/autostart/deactivate-caps-lock.desktop
[Desktop Entry]
Name=Deactivate CAPS LOCK
Exec=xmodmap /home/leandrost/.alt-keyboard
Terminal=true
Type=Application
StartupNotify=false
